# ha-mcp-bridge/rootfs/etc/s6-overlay/s6-rc.d/mcp-server/run
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Start MCP Server service - SIMPLIFIED VERSION
# ==============================================================================

bashio::log.info "Starting HA MCP Bridge Add-on..."
bashio::log.info "Version: 3.0.0 (API Proxy)"

# Load configuration
ADMIN_PASSWORD="$(bashio::config 'admin_password')"
NABU_CASA_URL="$(bashio::config 'nabu_casa_url')"

# Use fixed port 3003 as defined in config.yaml
export PORT="3003"
export LOG_LEVEL="info"

# Basic OAuth configuration
export OAUTH_ENABLED="true"
export ADMIN_USERNAME="admin"
export ADMIN_PASSWORD="${ADMIN_PASSWORD}"

# Home Assistant integration using supervisor token
export HA_URL="http://supervisor/core"
export HA_TOKEN="${SUPERVISOR_TOKEN}"

# CORS origins for Claude.ai
export CORS_ORIGIN="https://claude.ai,https://app.claude.ai"

# API Proxy via Ingress (Nabu Casa compatible)
CLAUDE_URL="https://ha.right-api.com/15715349_ha-mcp-bridge-v2/ingress"

bashio::log.info ""
bashio::log.info "🎯 =================================="
bashio::log.info "🎯    CLAUDE.AI SETUP INSTRUCTIONS"
bashio::log.info "🎯 =================================="
bashio::log.info ""
bashio::log.info "📋 Copy this URL for Claude.ai:"
bashio::log.info "🔗 ${CLAUDE_URL}"
bashio::log.info ""
bashio::log.info "📋 HA API Proxy via Nabu Casa ingress"
bashio::log.info "✅ Real HA API integration with supervisor token"
bashio::log.info "🔧 No authentication barriers for Claude.ai"
bashio::log.info ""
bashio::log.info "🎯 =================================="
bashio::log.info ""

# Change to application directory
cd /app

# Start the HA API Proxy (bypassing ingress)
bashio::log.info "🚀 Starting HA API Proxy server..."
bashio::log.info "🔄 Using HA API with supervisor token"
bashio::log.info "🎯 Serving via Nabu Casa ingress"
exec node /usr/bin/ha-api-proxy.js