# ha-mcp-bridge/rootfs/etc/s6-overlay/s6-rc.d/mcp-server/run
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Start MCP Server service - SIMPLIFIED VERSION
# ==============================================================================

bashio::log.info "Starting HA MCP Bridge Add-on..."
bashio::log.info "Version: 3.1.0 (SSE MCP Compatible)"

# Load configuration
ADMIN_PASSWORD="$(bashio::config 'admin_password')"
NABU_CASA_URL="$(bashio::config 'nabu_casa_url')"

# Use fixed port 3003 as defined in config.yaml
export PORT="3003"
export LOG_LEVEL="info"

# Basic OAuth configuration
export OAUTH_ENABLED="true"
export ADMIN_USERNAME="admin"
export ADMIN_PASSWORD="${ADMIN_PASSWORD}"

# Home Assistant integration using supervisor token
export HA_URL="http://supervisor/core"
export HA_TOKEN="${SUPERVISOR_TOKEN}"

# CORS origins for Claude.ai
export CORS_ORIGIN="https://claude.ai,https://app.claude.ai"

# SSE MCP via Ingress (like official HA integration)
CLAUDE_URL="https://ha.right-api.com/hassio/ingress/15715349_ha-mcp-bridge-v2/mcp_server/sse"

bashio::log.info ""
bashio::log.info "🎯 =================================="
bashio::log.info "🎯    CLAUDE.AI SETUP INSTRUCTIONS"
bashio::log.info "🎯 =================================="
bashio::log.info ""
bashio::log.info "📋 Copy this URL for Claude.ai:"
bashio::log.info "🔗 ${CLAUDE_URL}"
bashio::log.info ""
bashio::log.info "📋 SSE MCP Server - official HA integration compatible"
bashio::log.info "🌊 Server-Sent Events transport for Claude.ai remote MCP"
bashio::log.info "🔧 Endpoint: /mcp_server/sse"
bashio::log.info ""
bashio::log.info "🎯 =================================="
bashio::log.info ""

# Change to application directory
cd /app

# Start the HA API Proxy (bypassing ingress)
bashio::log.info "🚀 Starting SSE MCP server..."
bashio::log.info "🌊 SSE transport like official HA MCP integration"
bashio::log.info "🎯 Endpoint: /mcp_server/sse for Claude.ai"
exec node /usr/bin/sse-mcp-server.js